# cookbooks/server/templates/default/foresee-service.conf.erb
description "start and stop the foresee server"
version "1.0"
author "Varokas Panusuwan"

start on startup
stop on shutdown

respawn
respawn limit 15 5

env PORT=<%= @port %>

script
  export HOME="/root"
  exec /usr/bin/node <%= @path %>/app.js >> /var/log/<%= @logfile_name %> 2>&1
end script